# packaging/gateon.service
[Unit]
Description=Gateon - API Gateway and Reverse Proxy
Documentation=https://github.com/gateon/gateon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/gateon
Restart=on-failure
RestartSec=5s

# Config: place global.json, routes.json, etc. in /etc/gateon
WorkingDirectory=/var/lib/gateon
Environment=GLOBAL_CONFIG_FILE=/etc/gateon/global.json

# Override with: systemctl edit gateon
# Environment=PORT=8080

# State and config directories (auto-created with correct permissions)
StateDirectory=gateon
ConfigurationDirectory=gateon

# eBPF/XDP offloading: loading and attaching the embedded XDP program needs the
# BPF capabilities and an unbounded locked-memory limit for BPF map allocation
# (the loader also calls rlimit.RemoveMemlock(), but this keeps older kernels
# happy). The service runs as root so these caps are already held; they are
# declared explicitly so XDP keeps working if the unit is later de-rooted, and
# so NoNewPrivileges does not strip them.
LimitMEMLOCK=infinity
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN CAP_SYS_RESOURCE CAP_PERFMON

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=-/etc/gateon -/var/lib/gateon

[Install]
WantedBy=multi-user.target
